# container_confs/etc/profile.d/aliases.sh
# User defined aliases
alias dir='ls -alh --color'
